# .config/systemd/user/kdeconnectd.service
[Unit]
Description=KDE Connect daemon

[Service]
Type=simple
ExecStart=/usr/libexec/kdeconnectd

[Install]
WantedBy=default.target
